# tether.service
[Unit]
Description=Android USB Tethering

[Service]
Type=simple
ExecStart=/usr/local/bin/tether
