# tether@.service
[Unit]
Description=Android USB Tethering

[Service]
User=vmed
ExecStart=/home/vmed/bin/tether %I
