# tether@.service
[Unit]
Description=Android USB Tethering

[Service]
User=cimv2
ExecStart=/home/cimv2/bin/tether %I
